# scripts/com.mateotelfer.careeraiexplorer.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.mateotelfer.careeraiexplorer.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/mateotelfer/Claude Code - Working Folder/career-ai-explorer/server/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/mateotelfer/Claude Code - Working Folder/career-ai-explorer/server</string>
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key><false/>
        <key>NetworkState</key><true/>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key><string>3001</string>
        <key>NODE_ENV</key><string>production</string>
        <key>PATH</key><string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>ANALYZE_PROVIDER</key><string>ollama</string>
        <key>OLLAMA_URL</key><string>http://localhost:11434</string>
        <key>OLLAMA_MODEL</key><string>llama3.2:3b</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/mateotelfer/Claude Code - Working Folder/career-ai-explorer/scripts/logs/server.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mateotelfer/Claude Code - Working Folder/career-ai-explorer/scripts/logs/server.err.log</string>
    <key>ThrottleInterval</key><integer>5</integer>
</dict>
</plist>
